# backup-cron.timer
[Unit]
Description=Network Backup System - Minute Timer
Requires=backup-cron.service

[Timer]
# Estrategia: Timer cada minuto + early-exit en scheduled_runner.py
# Ventajas:
# - Cambios de horario desde web aplican inmediatamente (sin daemon-reload)
# - Más flexible que regenerar OnCalendar dinámicamente
# - Locking previene ejecuciones simultáneas
#
# El scheduled_runner.py consulta la DB y hace early-exit si no corresponde

OnCalendar=*:*:00
AccuracySec=30s
Persistent=false

[Install]
WantedBy=timers.target
