# deploy/backup_manager.service
# =============================================================================
# Systemd Unit - Backup Manager (Gunicorn)
# =============================================================================
# Archivo: /etc/systemd/system/backup_manager.service
#
# Path: /opt/backup_system (path real del deployment)
# =============================================================================

[Unit]
Description=Backup Manager Gunicorn Application Server
Documentation=https://github.com/Peripotec/backup_system
After=network.target
Wants=network.target

[Service]
Type=notify
User=root
Group=root
RuntimeDirectory=backup_manager

# Working directory
WorkingDirectory=/opt/backup_system

# Environment
Environment="PATH=/opt/backup_system/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# =============================================================================
# SOCKET BIND ABSOLUTO - Crítico para evitar 502
# =============================================================================
ExecStart=/opt/backup_system/venv/bin/gunicorn \
    --bind unix:/opt/backup_system/backup_manager.sock \
    --workers 2 \
    --threads 4 \
    --timeout 120 \
    --graceful-timeout 30 \
    --keep-alive 5 \
    --access-logfile /var/log/backup_manager/access.log \
    --error-logfile /var/log/backup_manager/error.log \
    --capture-output \
    web_app:app

# Socket permissions para Nginx
ExecStartPost=/bin/chmod 666 /opt/backup_system/backup_manager.sock

# Cleanup socket on stop
ExecStopPost=/bin/rm -f /opt/backup_system/backup_manager.sock

# =============================================================================
# Restart Policy - Self-healing
# =============================================================================
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
