# deploy/cert-watcher.timer
# =============================================================================
# Systemd Timer - Cert Watcher
# =============================================================================
# Archivo: /etc/systemd/system/cert-watcher.timer
# =============================================================================

[Unit]
Description=Certificate Watcher Timer
Documentation=https://github.com/Peripotec/backup_system

[Timer]
# Ejecutar cada 6 horas
OnCalendar=*-*-* 00,06,12,18:00:00
# Randomizar hasta 30 minutos para evitar thundering herd
RandomizedDelaySec=1800
# Ejecutar inmediatamente si se perdió la última ejecución
Persistent=true

[Install]
WantedBy=timers.target
